# dotfiles/profile.d/brew
#!/bin/bash

alias b=brew
alias s="brew search"

i() {
  brew install "$@" || brew cask install "$@"
}

ri() {
  brew reinstall "$@" || brew cask reinstall "$@"
}
