# dotfiles/profile.d/direnv
#!/bin/bash

export DIRENV_LOG_FORMAT=""
! which direnv >&- || eval "$(direnv hook bash)"
