# dotfiles/profile.d/direnv
#!/bin/bash

! which direnv >&- || eval "$(direnv hook bash)"
