# dotfiles/profile.d/direnv
#!/bin/bash

export DIRENV_LOG_FORMAT=""
! command -v direnv >&- || eval "$(direnv hook bash)"
alias e=direnv
