# dotfiles/profile.d/dns
#!/bin/bash

dns_flush() {
  sudo dscacheutil -flushcache
  sudo killall -HUP mDNSResponder
}

dns_info() {
  scutil --dns
}
