# dotfiles/profile.d/docker
#!/bin/bash

alias d=docker
alias dm=docker-machine
alias dc=docker-compose
