# dotfiles/profile.d/git
#!/bin/bash -e

unpair() {
  if which git >&-
  then
    pushd "$HOME/.setup" >&-
      git checkout dotfiles/gitconfig
    popd >&-
  fi
}
