# dotfiles/profile.d/git
#!/bin/bash -e

alias g=git
# tig, git's best companion
alias t="tig --all"

unpair() {
  if which git >&-
  then
    (
      cd "$HOME/.setup" && git checkout dotfiles/gitconfig
    )
  fi
}
