# dotfiles/profile.d/go
#!/bin/bash

export PATH="$HOME/go/bin:$PATH"
