# dotfiles/profile.d/go
#!/bin/bash

if which go >&-
then
  export GOPATH="${HOME}/go"
  export PATH="${GOPATH}/bin:/usr/local/opt/go/libexec/bin:${PATH}"
fi
