# dotfiles/profile.d/kubernetes
#!/bin/bash

alias k="kubectl"
alias h="helm"
