# dotfiles/profile.d/make
#!/bin/bash

alias m=make
alias boss=make

# make is now GNU Make v4
export PATH="/usr/local/opt/make/libexec/gnubin:$PATH"
