# dotfiles/profile.d/make
#!/bin/bash

alias m=gmake
alias boss=gmake
