# dotfiles/profile.d/pw
#!/bin/bash

pw() {
  sf-pwgen -l 20 -c 10 "$@"
}
