# dotfiles/profile.d/rbenv
#!/bin/bash

if which rbenv >&-
then
  export RBENV_ROOT=/usr/local/var/rbenv
  eval "$(rbenv init -)"
fi
