# dotfiles/profile.d/ruby
#!/bin/bash

alias bx="bundle exec"

alias bl="bundle list"

alias bi="bundle install"
alias bu="bundle update"

alias bout="bundle outdated"
