# dotfiles/profile.d/snatch
#!/bin/bash

! [ -d $HOME/Desktop/snatch ] || eval "$($HOME/Desktop/snatch/bin/snatch init -)"
