# dotfiles/profile.d/ssl
#!/bin/bash

function verify_ssl() {
  local fqdn
  fqdn="$1"

  openssl s_client -showcerts -servername "$fqdn" -connect "$fqdn:443"
}
