# dotfiles/profile.d/ssl
#!/bin/bash

function verify_ssl() {
  local url
  url="$1"

  openssl s_client -showcerts -connect "$url:443"
}
