# dotfiles/profile.d/terminal
#!/bin/bash

# enable ls coloring
export CLICOLOR=1

# List info about this Mac
# https://github.com/dylanaraps/neofetch
alias mac=neofetch
