# dotfiles/profile.d/vim
#!/bin/bash

export EDITOR=vim

alias v=vim
